# crontab/run_spider.sh
#!/usr/bin/env bash

python apps/spiders/baic_base/setup_service.py &