# crontab/run_spider.sh
#!/usr/bin/env bash

python apps/run_spider.py